# 01-storage/configs/etc/cron.weekly/btrfs-maintenance
#!/bin/bash
btrfs scrub start /mnt/data
btrfs balance start -dusage=85 /mnt/data
